# home/Library/LaunchAgents/homebrew-doctor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
        <string>localhost.homebrew-doctor</string>
    <key>ProcessType</key>
        <string>Background</string>
    <key>ProgramArguments</key>
        <array>
            <string>/bin/sh</string>
            <string>-c</string>
            <string>/usr/local/bin/brew update &amp;&amp; /usr/local/bin/brew doctor; /usr/local/bin/terminal-notifier -title 'Homebrew updated!!' -message 'see /tmp/localhost.homebrew-doctor*' </string>
        </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/localhost.homebrew-doctor.stderr</string>
    <key>StandardOutPath</key>
    <string>/tmp/localhost.homebrew-doctor.stdout</string>
    <key>StartCalendarInterval</key>
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>1</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
</dict>
</plist>
